# metrics/sender/cron/crontab
*/1 * * * * cd /mydata/firm/metrics/sender && sh /mydata/firm/metrics/sender/cron/bin/sender.sh 2>&1 | /usr/bin/logger -t sender
